# rag-systemd/rag.service
[Unit]
Description=RAG Test Service - Document Monitoring
Documentation=file:///home/liudecheng/systemd-rag/README.md
After=network.target ollama.service
Wants=ollama.service
Requires=ollama.service

[Service]
Type=simple
User=liudecheng
Group=liudecheng
ExecStart=/home/liudecheng/.local/share/ulauncher/extensions/com.github.497672776.ulauncher-rag-retriever/rag-systemd/start_service.sh

# 工作目录
WorkingDirectory=/home/liudecheng/.local/share/ulauncher/extensions/com.github.497672776.ulauncher-rag-retriever/rag-systemd

# 环境变量
Environment=PYTHONPATH=.
Environment=PYTHONUNBUFFERED=1

# 重启策略
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30

# 资源限制
MemoryLimit=1G
CPUQuota=100%
TasksMax=50

# 日志配置
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target